# .config/systemd/user/mail-sync.timer
[Unit]
Description = Mail sync timer

[Timer]
OnBootSec = 10m
OnUnitActiveSec = 30m
Unit = mail-sync.service

[Install]
WantedBy = timers.target
