# .config/systemd/user/mail-sync.service
[Unit]
Description = Mail sync service

[Service]
Type=oneshot
ExecStart = /home/cweagans/bin/mail-sync
