# .config/systemd/user/msmtp-queue.timer
[Unit]
Description=Timer to flush msmtp queue

[Timer]
OnUnitInactiveSec=10m

[Install]
WantedBy=timers.target
