# .config/systemd/user/todoist-sync.service
[Unit]
Description = Todoist sync service

[Service]
Type=oneshot
ExecStart = /home/cweagans/bin/todoist-sync
