# .config/systemd/user/todoist-sync.service
[Unit]
Description = Todoist sync service

[Service]
Type=oneshot
ExecStart = /usr/bin/todoist sync
